# systemd/cloudflare-ddns.timer
[Unit]
Description=Update CloudFlare DNS for every 15 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target